# run_shells/infer/multi_turn_conversation/web/llama_gpt35sex.conf
[program: llama_gpt35sex]
command = gunicorn -c gunicorn_conf_gpt35sex.py model_infer_flask_server_gpt35sex:app  --error-logfile /mnt/cephfs/hjh/logs/llama_gpt35sex/supervisor_err.log
directory = /home/huangjiahong.dracu/hjh/pycharm_projects/nlp/stanford_alpaca/run_shells/infer/multi_turn_conversation/web
user = huangjiahong.dracu ;
autostart=true ;
autorestart=true ;
redirect_stderr=true ;
logfile_backups=10 ;
stdout_logfile_backups=10 ;
logfile_maxbytes=25MB ;
stderr_logfile=/mnt/cephfs/hjh/logs/llama_gpt35sex/supervisor_err.log ;
stdout_logfile=/mnt/cephfs/hjh/logs/llama_gpt35sex/supervisor_info.log ;
